# ansible/roles/common/templates/schedule-passes.service.j2
[Unit]
Description=Schedule passes

[Service]
ExecStart={{ noaa_home }}/scripts/schedule.sh -t
User=pi
